# deploy/fitness-backend.service
[Unit]
Description=Fitness Management System Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fitness-system/backend
Environment="SPRING_PROFILES_ACTIVE=prod"
Environment="JWT_SECRET=please-change-this-secret-key"

# Java启动命令
ExecStart=/usr/bin/java -jar \
    -Xms256m \
    -Xmx512m \
    -Dspring.profiles.active=prod \
    -Dspring.config.location=/opt/fitness-system/backend/application-prod.yml \
    /opt/fitness-system/backend/fitness-system-1.0.0.jar

# 重启策略
Restart=always
RestartSec=10

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fitness-backend

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
